# deploy/standby-tool.service
[Unit]
Description=Standby-Tool Backend (FastAPI, serviert auch das Frontend)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# User + Pfade ggf. an den Pi anpassen (Standard: Benutzer "pi", Repo unter Home).
User=pi
WorkingDirectory=/home/pi/SMART-Tool
# Site-spezifische Env (Shelly-IP etc.). "-" = optional, Start auch ohne Datei.
EnvironmentFile=-/home/pi/SMART-Tool/deploy/standby-tool.env
ExecStart=/home/pi/SMART-Tool/.venv/bin/uvicorn app.main:app --app-dir backend --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
